# guardian/launchd/io.moltbot.hardened.guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.moltbot.hardened.guardian</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Developer/CommandLineTools/usr/bin/python3</string>
        <string>/usr/local/lib/moltbot-hardened/guardian.py</string>
        <string>--interval</string>
        <string>30</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/moltbot-hardened/guardian.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/moltbot-hardened/guardian.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MBH_GUARDIAN_LOG</key>
        <string>/usr/local/var/log/moltbot-hardened/guardian.log</string>
        <key>MBH_CLI</key>
        <string>/usr/local/bin/moltbot-hardened</string>
    </dict>
</dict>
</plist>
